# systemd/fath-bot.service
[Unit]
Description=FATH Unified Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fath-bot
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /opt/fath-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

